# service/log/run
#!/bin/sh

exec /usr/bin/chpst -u logger /usr/bin/svlogd -tt main
